# connman/src/connman.service.in
[Unit]
Description=Connection service
Requires=@CONNMAN_SERVICE_REQUIRES@
After=@CONNMAN_SERVICE_AFTER@ network-pre.target
Before=network.target
Wants=network.target

[Service]
Type=notify
Restart=always
EnvironmentFile=-/etc/sysconfig/connman
EnvironmentFile=-/var/lib/environment/connman/*.conf
ExecStart=@sbindir@/connmand -n -W nl80211 --nobacktrace --systemd --noplugin=wifi $SYSCONF_ARGS $CONNMAN_ARGS
ExecReload=/usr/bin/dbus-send --system --print-reply --type=method_call --dest=net.connman / net.connman.Firewall.Reload
StandardOutput=null

[Install]
WantedBy=multi-user.target
