# connman/src/connman.service.in
[Unit]
Description=Connection service
Requires=dbus.service
After=dbus.service
Before=remote-fs.target

[Service]
Type=notify
Restart=on-failure
EnvironmentFile=-/etc/tracing/connman
ExecStart=@prefix@/sbin/connmand -n --systemd $TRACING
StandardOutput=null

[Install]
WantedBy=multi-user.target
