# connman/src/connman.service.in
[Unit]
Description=Connection service
Requires=@CONNMAN_SERVICE_REQUIRES@
After=@CONNMAN_SERVICE_AFTER@
Before=remote-fs.target

[Service]
Type=notify
Restart=always
EnvironmentFile=-/etc/tracing/connman/connman.tracing
ExecStart=@prefix@/sbin/connmand -n --nobacktrace --systemd $TRACING
StandardOutput=null

[Install]
WantedBy=multi-user.target
