# connman/src/connman.service.in
[Unit]
Description=Connection service
After=syslog.target

[Service]
Type=dbus
BusName=net.connman
Restart=on-failure
EnvironmentFile=-/etc/tracing/connman
ExecStart=@prefix@/sbin/connmand -n $TRACING
StandardOutput=null
Restart=always

[Install]
WantedBy=multi-user.target
