# connman/src/connman.service.in
[Unit]
Description=Connection service
Requires=dbus.service
After=dbus.service

[Service]
Type=notify
Restart=on-failure
EnvironmentFile=-/etc/tracing/connman
ExecStart=@prefix@/sbin/connmand -n --systemd $TRACING
StandardOutput=null
Restart=always

[Install]
WantedBy=multi-user.target
